# systemd/mvwifi-auto-resume.service
[Unit]
Description=MV WiFi Auto Resume Check
After=suspend.target hibernate.target hybrid-sleep.target NetworkManager.service

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 5
ExecStart=%h/.local/bin/mvwifi-auto --once
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=suspend.target hibernate.target hybrid-sleep.target
